# supervisor/conf.d/laravel-vue-admin.conf
[program:laravel-vue-admin]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work redis --timeout=360 --sleep=3 --tries=3
autostart=true
autorestart=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stopwaitsecs=3600